# packages/kafka/zookeeper.service
[Unit]
Description=zookeeper service
After=network.target

[Service]
ExecStart=/opt/kafka/bin/zookeeper-server-start.sh /opt/kafka/config/zookeeper.properties
Restart=on-failure
User=kafka
Group=kafka

[Install]
WantedBy=multi-user.target
